# conf/systemd.service
[Service]
Environment="JELLYFIN_PublishedServerUrl=https://__DOMAIN____PATH__"
